# MSEE/supervisor/arp_mock_server.conf
[program:arp_mock_server]
command=/usr/sbin/arp_mock_server
priority=3
autostart=true
autorestart=false
stdout_logfile=/tmp/arp_mock_server.out.log
stderr_logfile=/tmp/arp_mock_server.err.log

