# MSEE/supervisor/arp_responder_msee.conf
[program:arp_responder]
command=/usr/sbin/arpresponder_msee
priority=3
autostart=true
autorestart=false
stdout_logfile=/tmp/arp_responder.out.log
stderr_logfile=/tmp/arp_responder.err.log

