# MSEE/arp_responder/tests/build/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
process_name=sshd
stdout_logfile=/tmp/sshd.out.log
stderr_logfile=/tmp/sshd.err.log
redirect_stderr=false
autostart=true
autorestart=true
startsecs=1
numprocs=1

[program:arpresponder_bm]
command=/usr/sbin/arpresponder_bm
process_name=arpresponder_bm
stdout_logfile=/tmp/arpresponder_bm.out.log
stderr_logfile=/tmp/arpresponder_bm.err.log
redirect_stderr=false
autostart=false
autorestart=false
startsecs=1
numprocs=1

[program:arpresponder_msee]
command=/usr/sbin/arpresponder_msee
process_name=arpresponder_msee
stdout_logfile=/tmp/arpresponder_msee.out.log
stderr_logfile=/tmp/arpresponder_msee.err.log
redirect_stderr=false
autostart=false
autorestart=false
startsecs=1
numprocs=1

