# supervisor/redis.conf
[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf
priority=1
autostart=true
autorestart=false
stdout_logfile=/tmp/redis-out.log
stderr_logfile=/tmp/redis-err.log
