# supervisor/rest_api_test.conf
[program:rest-api]
command=/usr/sbin/go-server-server.test -test.coverprofile=/coverage.cov -systemtest=true -enablehttps=true -clientcert=/usr/sbin/cert/client/selfsigned.crt -servercert=/usr/sbin/cert/server/selfsigned.crt -serverkey=/usr/sbin/cert/server/selfsigned.key -localapitestdocker=true -clientcertcommonname=test.client.restapi.sonic,*.example.sonic,*test.sonic,*. -loglevel trace
priority=1
autostart=false
autorestart=false
stdout_logfile=/tmp/rest-api.out.log
stderr_logfile=/tmp/rest-api.err.log

[program:start.sh]
command=/usr/bin/start.sh
priority=1
autostart=true
autorestart=false
stdout_logfile=/tmp/start.out.log
stderr_logfile=/tmp/start.err.log

[program:CreateMockPort.sh]
command=/usr/bin/CreateMockPort.sh
priority=1
autostart=true
autorestart=false
stdout_logfile=/tmp/CreateMockPort.out.log
stderr_logfile=/tmp/CreateMockPort.err.log

[program:ServerKill]
command=/usr/bin/server_kill.py
priority=1
autostart=true
autorestart=false
